# conf/nextcloud.cron
*/15  *  *  *  * __USER__ /usr/bin/php__PHPVERSION__ -f __FINALPATH__/cron.php
