# conf/nextcloud.cron
*/15  *  *  *  * #USER# /usr/bin/php -f #DESTDIR#/cron.php
