# conf/nextcloud.cron
*/15  *  *  *  * __USER__ /usr/bin/php__YNH_PHP_VERSION__ -f __DESTDIR__/cron.php
