# conf/nextcloud.cron
*/15  *  *  *  * __USER__ /usr/bin/php -f __DESTDIR__/cron.php
