# conf/nextcloud.cron
*/15  *  *  *  * __APP__ /usr/bin/php__PHPVERSION__ -f __FINAL_PATH__/cron.php
