# conf/nextcloud.cron
*/15  *  *  *  * __APP__ /usr/bin/php__PHPVERSION__ --define apc.enable_cli=1 -f __FINAL_PATH__/cron.php
